# debian/postinst
#!/bin/sh
set -e

if [ "$1" = "configure" ]; then
    systemctl daemon-reload
    systemctl enable armbian-web-config.service || true
fi

#DEBHELPER#
exit 0
