# debian/prerm
#!/bin/sh
set -e

if [ "$1" = "remove" ]; then
    systemctl disable --now armbian-web-config.service || true
fi

#DEBHELPER#
exit 0
